# systemd/mntdata.service
[Unit]
Description=Mount CFSR data

[Service]
ExecStart=/root/mntdata

[Install]
WantedBy=multi-user.target

