# systemd/sailing.service
[Unit]
Description=Sailing

[Service]
ExecStart=/root/sailing/startsailing

[Install]
WantedBy=multi-user.target

